# private_dot_config/systemd/user/nixos-vscode.path.tmpl
[Unit]
Description=VSCode server directory

[Path]
MakeDirectory=true
PathChanged={{ .chezmoi.homeDir }}/.vscode-server/bin

[Install]
WantedBy=paths.target
